# cli.service
[Unit]
Description=Bermuda cli

[Service]
ExecStart=python /home/pi/bermuda/cli.py
