# compass.service
[Unit]
Description=Bermuda kompas

[Service]
ExecStart=python /home/pi/bermuda/compass.py
