# radar.service
[Unit]
Description=Bermuda radar

[Service]
ExecStart=python /home/pi/bermuda/radar.py
